# flightdeck.service.example
[Unit]
Description=Flightdeck
After=network.target

[Service]
User=flightdeck
WorkingDirectory=/home/flightdeck/flightdeck
EnvironmentFile=-/home/flightdeck/flightdeck/.env
ExecStart=/home/flightdeck/flightdeck/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
